# crontab
0 6 * * *     python daily-reset.py
0 6 * * TUE   python weekly-reset.py
0 6 1 * *     python monthly-reset.py
