# crontab
 0 6 * * *    python /home/virginitybot/daily-reset.py
 0 7 * * TUE  python /home/virginitybot/weekly-reset.py
30 7 1 * *    python /home/virginitybot/monthly-reset.py
